# dist/linux/runit/saffroncoind-testnet/run
#!/bin/sh

USER=csfrd
USER_HOME=/home/csfr

exec chpst -u ${USER} saffroncoind -daemon=0 -testnet -datadir=${USER_HOME}/.saffroncoin-testnet