# dist/linux/runit/insight/run
#!/bin/sh

USER=csfrd
USER_HOME=/home/csfr
chown -R $USER $USER_HOME/insight-api
export INSIGHT_NETWORK=livenet
export INSIGHT_DB=$USER_HOME/insight-api/db
export BITCOIND_DATADIR=$USER_HOME/.saffroncoin/
export BITCOIND_USER=`cat $USER_HOME/.saffroncoin/saffroncoin.conf | sed -n 's/.*rpcuser=\([^ \n]*\).*/\1/p'`
export BITCOIND_PASS=`cat $USER_HOME/.saffroncoin/saffroncoin.conf | sed -n 's/.*rpcpassword=\([^ \n]*\).*/\1/p'`
#BITCOIND_HOST -- use the default (localhost)
#BITCOIND_PORT -- use the default (19710)
#BITCOIND_P2P_PORT -- use the default (19717)

#sv start bitcoind || exit 1
cd ${USER_HOME}/insight-api
exec 2>&1
echo "Waiting 60 seconds before starting insight, to allow saffroncoind to fully initialize..."
sleep 60
exec chpst -u ${USER} /usr/bin/node ${USER_HOME}/insight-api/insight.js
