# dist/linux/runit/counterpartyd-testnet/run
#!/bin/sh

USER=csfrd
export USER_HOME=/home/csfr

exec su -s /bin/bash -c '/usr/local/bin/csfrd --data-dir=${USER_HOME}/.config/csfrd-testnet --testnet server' ${USER}
